# saiman-signal.service
[Unit]
Description=Saiman Signal Bot
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/saiman-signal
EnvironmentFile=/home/ubuntu/saiman-signal/.env
ExecStart=/home/ubuntu/.local/bin/uv run python -m saiman_signal
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
